# config/supervisor.conf
[program:imagelab]
command=/var/www/imagelab-backend/venv/bin/gunicorn -c config/gunicorn.conf.py app:app
directory=/var/www/imagelab-backend
user=www-data
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopwaitsecs=30
environment=
    PATH="/var/www/imagelab-backend/venv/bin:%(ENV_PATH)s"

stdout_logfile=/var/www/imagelab-backend/logs/supervisor_stdout.log
stderr_logfile=/var/www/imagelab-backend/logs/supervisor_stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

[program:celery-worker]
command=/var/www/imagelab-backend/venv/bin/celery -A app.celery worker --loglevel=info
directory=/var/www/imagelab-backend
user=www-data
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopwaitsecs=30

stdout_logfile=/var/www/imagelab-backend/logs/celery_stdout.log
stderr_logfile=/var/www/imagelab-backend/logs/celery_stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

[program:celery-beat]
command=/var/www/imagelab-backend/venv/bin/celery -A app.celery beat --loglevel=info
directory=/var/www/imagelab-backend
user=www-data
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopwaitsecs=30

stdout_logfile=/var/www/imagelab-backend/logs/celery_beat_stdout.log
stderr_logfile=/var/www/imagelab-backend/logs/celery_beat_stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

[group:imagelab-group]
programs=imagelab,celery-worker,celery-beat
priority=999
